# crontab
# Caps Edge data refresh schedule
# Runs once daily at 6:00 AM ET

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app
DATA_DIR=/app/data

# 6:00 AM ET (06:00)
0 6 * * * root cd /app && python scripts/refresh.py >> /var/log/caps-edge-refresh.log 2>&1

# Empty line required at end of crontab
